# backups.crontab
15 00 * * * /home/backups/pg_backup_rotated.sh
