# databasebackup.timer
[Unit]
Description=Database Backups

[Timer]
OnBootSec=5m
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
